# deploy/volume-service/synapse-mastervol.service
[Unit]
Description=Synapse master volume (control daemon + jack_mix_box gain stage)
After=modep-mod-host.service modep-mod-ui.service jack.service
BindsTo=jack.service
Wants=modep-mod-ui.service

[Service]
# volumectl.py is the tracked main process: it spawns/owns jack_mix_box
# ("synapsevol"), listens for raw volume CCs on its "control in" port (bus tap
# + app slider), applies the taper and drives the mixer over a private CC7
# link, and echoes applied state on "state out". The mixer keeps passing audio
# at its last gain even if the app/UI dies. RT scheduling is inherited by the
# spawned jack_mix_box. Runs as `modep` (in group `jack`) to reach the
# promiscuous server. Needs the JACK-Client python package (see README).
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=80
LimitRTPRIO=99
LimitMEMLOCK=infinity
User=modep
Group=modep
Environment=JACK_PROMISCUOUS_SERVER=jack

ExecStart=/usr/bin/python3 /usr/local/bin/synapse-volumectl.py
# insert the mixer into mod-monitor -> system:playback once its ports are up
ExecStartPost=/usr/local/bin/synapse-mastervol-rewire.sh synapsevol
# on stop/restart, restore the direct path so audio survives
ExecStopPost=/usr/local/bin/synapse-mastervol-revert.sh synapsevol
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
